# ops/kiosk/mini-pc/snapshot/home/kiosk-j/dot-config/systemd/user/kiosk-url-guard.service
[Unit]
Description=Kiosk external URL guard
After=kiosk-chrome.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/kiosk-j/.local/bin/kiosk-url-guard.py
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
